# PREPARE/files/mogilefs/profile.d/mogilefs.sh

export PATH="/opt/MogileFS/bin:$PATH"
export PERL5LIB="/opt/MogileFS/lib/perl5"